# rdp_monitor/task_scheduler/RDP_Stop.xml
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>RDP Stop Log Handler for EventID 24 and 23</Description>
  </RegistrationInfo>

  <Triggers>
    <EventTrigger>
      <Enabled>true</Enabled>
      <Subscription><![CDATA[
<QueryList>
  <Query Id="0" Path="Microsoft-Windows-TerminalServices-LocalSessionManager/Operational">
    <Select>*[System[(EventID=24 or EventID=23)]]</Select>
  </Query>
</QueryList>
      ]]></Subscription>
    </EventTrigger>
  </Triggers>

  <Principals>
    <Principal id="Author">
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>

  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
  </Settings>

  <Actions Context="Author">
    <Exec>
      <Command>C:\Windows\System32\WindowsPowerShell\v1.0\powershell.exe</Command>
      <Arguments>-ExecutionPolicy Bypass -NoProfile -File "C:\rdp-video-pjt\rdp_monitor\scripts\rdp_stop_log.ps1"</Arguments>
      <WorkingDirectory>C:\rdp_monitor\scripts</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
